# deploy/sommelier_frontend.service
[Unit]
Description=Sommelier Frontend (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/sommelier_frontend
ExecStart=next start -p 3000
Restart=always
RestartSec=5
KillMode=control-group
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin:/srv/sommelier_frontend/node_modules/.bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
